# cron/crontab
*       *       *       *       *       /path/to/usrmgt.sh
*       *       *       *       *       /path/to/ttenv.sh
*       *       *       *       *       /path/to/usrsyn.sh